# deploy/devspace-bot.service
[Unit]
Description=DEVSPACE Telegram Bot
After=network.target mysql.service

[Service]
Type=simple
WorkingDirectory=/opt/devspace-bot
EnvironmentFile=/opt/devspace-bot/.env
ExecStart=/opt/devspace-bot/.venv/bin/python -m devspace_bot
Restart=always
RestartSec=5
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target

